# etc/systemd/system/set-ryzenadj-tweaks.path
[Path]
PathChanged="$HOME"/.local/bin/allowadj.txt

[Install]
WantedBy=multi-user.target
