# etc/systemd/system/set-ryzenadj-tweaks.path.template
[Path]
PathChanged=__HOME__/.local/share/ryzen_uv/allowadj.txt

[Install]
WantedBy=multi-user.target
